# supervisor_anistream.conf
[program:anistream]
command=/home/user/anistream/venv/bin/gunicorn src.wsgi:application --bind 0.0.0.0:8000 --workers 3 --timeout 60 --access-logfile /var/log/anistream/gunicorn_access.log --error-logfile /var/log/anistream/gunicorn_error.log
directory=/home/user/anistream
user=www-data
autostart=true
autorestart=true
stderr_logfile=/var/log/anistream/supervisor_error.log
stdout_logfile=/var/log/anistream/supervisor_stdout.log

# Restart sozlamalari
startsecs=10
stopasgroup=true
killasgroup=true

# Environment variables (agar kerak bo'lsa)
environment=PATH="/home/user/anistream/venv/bin",HOME="/home/user/anistream"

